# provisioning/templates/whatson-ingest.timer
[Unit]
Description=Run whatson ingest repeatedly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
